# debian/eat-my-sms.target
[Unit]
Description=Eat My SMS - All GSM Modems (Auto-detected)
Documentation=https://github.com/pescheckit/eat-my-sms
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/eat-my-sms-start-all

[Install]
WantedBy=multi-user.target
